# https_dns_proxy.service.in
[Unit]
Description=https-dns-proxy - A light-weight DNS-HTTPS, non-caching translation proxy
Requires=network.target
Wants=nss-lookup.target
Before=nss-lookup.target
After=network.target

[Service]
Type=simple
Restart=on-failure
ExecStart=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_BINDIR}/https_dns_proxy \
  -u nobody -g nogroup -v -v ${SERVICE_EXTRA_OPTIONS}
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
